# nginx-gunicorn-flask-web/gunicorn.conf
[program:gunicorn]
command=gunicorn run:app -b localhost:5000
directory=/deploy/app
priority=999
stderr_logfile=/var/log/candybudget.err.log
stdout_logfile=/var/log/candybudget.out.log