# nginx-gunicorn-flask-web/mongod.conf
[program:mongod]
command=/usr/bin/mongod
priority=1
stderr_logfile=/var/log/mongodb.err.log
stdout_logfile=/var/log/mongodb.out.log
