# nginx-gunicorn-flask-web/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
priority=555
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log
